# Synced/syncRunner/workingStartup.service
Unit]
Description=Startup

[Service]
ExecStart=/home/petr/Documents/LocalShared/syncRunner/test.sh

[Install]
WantedBy=multi-user.target
